# deploy/avifeed.service
[Unit]
Description=avifeed — random AVIF image server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jhll1124
WorkingDirectory=/server/avifeed
# Run node directly so systemd's MainPID is the actual server process and
# SIGTERM is delivered to it without any wrapper in between. `pnpm build`
# produces dist/index.js; package.json sets "main": "dist/index.js".
ExecStart=/usr/local/bin/node dist/index.js
Environment=NODE_ENV=production
Environment=RIS_LOG_LEVEL=info
Restart=on-failure
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
